# roles/minecraft-ic2/templates/clean_logs.timer.j2
[Unit]
Description=Run Minecraft log cleanup daily

[Timer]
OnCalendar=*-*-* {{ log_cleanup_time }}
Persistent=true

[Install]
WantedBy=timers.target
