# roles/minecraft-ic2/templates/clean_logs.service.j2
[Unit]
Description=Clean old Minecraft logs

[Service]
Type=oneshot
ExecStart={{ minecraft_home }}/scripts/clean_logs.sh
